# daemon/plasma-powerdevil.service.in
[Unit]
Description=Powerdevil
PartOf=graphical-session.target
After=plasma-core.target

[Service]
ExecStart=@KDE_INSTALL_FULL_LIBEXECDIR@/org_kde_powerdevil
Type=dbus
BusName=org.kde.Solid.PowerManagement
TimeoutSec=5sec
Slice=background.slice
Restart=on-failure
